# deploy/systemd/openclaw-vnc-tunnel.service
[Unit]
Description=OpenClaw VM VNC tunnel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/ssh -o BatchMode=yes -o ConnectTimeout=10 -o TCPKeepAlive=yes -o ServerAliveInterval=15 -o ServerAliveCountMax=2 -o ExitOnForwardFailure=yes -N -L 127.0.0.1:5901:127.0.0.1:5901 openclaw-vm
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
