# deploy/hermes_openclaw_compat/hermes-openclaw-compat.service
[Unit]
Description=Hermes OpenClaw Compatibility Layer
After=network-online.target hermes-api-server.service
Requires=hermes-api-server.service

[Service]
Type=simple
EnvironmentFile=%h/.config/clawcontrol-hermes-compat.env
ExecStart=%h/.hermes/hermes-agent/venv/bin/python %h/.local/share/clawcontrol-hermes/hermes_openclaw_compat.py serve
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
