# deploy/hermes_openclaw_compat/hermes-api-server.service
[Unit]
Description=Hermes API Server
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/clawcontrol-hermes.env
WorkingDirectory=%h/.hermes/hermes-agent
ExecStart=%h/.local/bin/hermes gateway run
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
